# Daemons/Supervisor/confs/template.conf
[program:{daemonName}]
command={binDir}/console bozos:daemons:start {daemonName}
autostart=true
autorestart=true
stopsignal=TERM
user=www-data
stdout_logfile={logFile}