# Daemons/Supervisor/confs/template.conf
[program:{daemonName}]
command={binDir}/console bozos:daemons:start {daemonName} --no-debug
autostart={autostart}
autorestart=true
stopsignal=TERM
user=www-data
stdout_logfile={logFile}