# Daemons/Supervisor/confs/template.conf
[program:{daemonName}]
command={binDir}/console daemons:start {daemonName}
autostart=false
autorestart=true
stopsignal=TERM
user=www-data
stdout_logfile={logFile}